# src/daemon/data/ddccontrol.service.in
[Unit]
Description=Control monitor parameters, like brightness, contrast, and other...

[Service]
Type=dbus
BusName=ddccontrol.DDCControl
ExecStart=@bindir@/ddccontrol_service
User=root
