# src/daemon/data/ddccontrol.service.in
[Unit]
Description=Control monitor parameters, like brightness, contrast, and other...
Documentation=man:ddccontrol(1)

[Service]
Type=dbus
BusName=ddccontrol.DDCControl
ExecStart=@servicedir@/ddccontrol_service
User=root
